# crontab
# Run on friday at 00:15 and upon start
#15 0 * * 5 /renew.sh > /proc/1/fd/1 2>/proc/1/fd/2
0 16 * * * /renew.sh > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
